# pybackup.timer
[Unit]
Description=Timer for pybackup

[Timer]
OnCalendar=*-*-* 2:00:00
Persistent=true

[Install]
WantedBy=timers.target
