# pybackup_timer.service
[Unit]
Description=Timer for pybackup.

[Timer]
OnCalendar=*-*-* 22:09:00
